# services/projects/polymarket-insider-bot.service
[Unit]
Description=Polymarket Insider Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wolf
Group=wolf
WorkingDirectory=/home/wolf/projects/polymarket-insider-bot
ExecStart=/home/wolf/.local/bin/mise run daemon
Restart=on-failure
RestartSec=30
Environment=HOME=/home/wolf
# Store databases/data on HDD to avoid filling the SSD
Environment=DATABASE_URL=sqlite+aiosqlite:///mnt/storage/projects/polymarket-insider-bot/polyinsider.db

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/wolf/projects/polymarket-insider-bot
ReadWritePaths=/mnt/storage/projects

[Install]
WantedBy=multi-user.target
